# crontab
#updates the longest_uptime file every minute
*/1 * * * * /home/<user>/uptimer.sh /home/<user>/longest_uptime

